# docs/systemd/parlaygorilla.service
# systemd unit for Parlay Gorilla (Oracle VM).
# Install: sudo cp docs/systemd/parlaygorilla.service /etc/systemd/system/
#         sudo systemctl daemon-reload
#         sudo systemctl enable parlaygorilla
# Replace /opt/parlaygorilla with your deploy path if different.

[Unit]
Description=Parlay Gorilla API and Scheduler (Docker Compose)
After=docker.service network-online.target
Wants=docker.service network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/parlaygorilla
ExecStart=/usr/bin/docker compose -f /opt/parlaygorilla/docker-compose.prod.yml up -d --build --remove-orphans
ExecStartPost=/usr/bin/docker compose -f /opt/parlaygorilla/docker-compose.prod.yml ps
# Fail start if nginx isn't responding on localhost:80 within ~30s.
ExecStartPost=/bin/bash -lc 'set -e; for i in {1..30}; do if curl -fsS --max-time 2 http://127.0.0.1/health >/dev/null; then exit 0; fi; sleep 1; done; echo "[systemd] ERROR: nginx /health not responding on :80" >&2; exit 1'
# Fail start if API liveness isn't responding on localhost:8000 within ~90s (migrations can take time).
ExecStartPost=/bin/bash -lc 'set -e; for i in {1..90}; do if curl -fsS --max-time 2 http://127.0.0.1:8000/healthz >/dev/null; then exit 0; fi; sleep 1; done; echo "[systemd] ERROR: api /healthz not responding on :8000" >&2; exit 1'
ExecStop=/usr/bin/docker compose -f /opt/parlaygorilla/docker-compose.prod.yml down
TimeoutStartSec=600
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
