# docs/systemd/parlaygorilla-watchdog.service
# Oneshot service run by parlaygorilla-watchdog.timer.
# If origin /health fails, restarts the parlaygorilla stack.

[Unit]
Description=Parlay Gorilla health watchdog (check origin /health)
After=network-online.target parlaygorilla.service
Wants=network-online.target

[Service]
Type=oneshot
Environment=APP_DIR=/opt/parlaygorilla
Environment=HEALTH_URL=http://127.0.0.1/health
Environment=MAX_WAIT=5
ExecStart=/opt/parlaygorilla/scripts/health-watchdog.sh
# Do not let watchdog restart loop; timer rate-limits.
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
